# systemd/parent-bot.service
[Unit]
Description=PocketBot - Parent Bot
After=network.target

[Service]
WorkingDirectory=/opt/pocketbot
Environment="PYTHONPATH=/opt/pocketbot"
ExecStart=/usr/bin/python -m app.bots.parent.main
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
